# backend/docker/crontab
* * * * * /var/www/backend/bin/console app:send-notifications
